# .config/systemd/user/notes-backup.path
[Unit]
Description=Watches notes directory for sync to NAS

[Path]
PathChanged=%h/notes

[Install]
WantedBy=default.target
